# cron/satis.cron
* * * * * root bash /satis/bin/refresh.sh >> /var/log/cron.log
